# infra/systemd/gpt-dev-registry-rules.service
[Unit]
Description=Hetzner Dev Workspace - Registry Network and iptables Rules
After=docker.service
Requires=docker.service
BindsTo=docker.service
PartOf=docker.service
Before=gpt-dev-gateway.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /opt/hetzner-dev-workspace/infra/scripts/setup-registry-network.sh

[Install]
WantedBy=multi-user.target docker.service
