# infra/systemd/gpt-dev-gateway.service
[Unit]
Description=Hetzner Dev Workspace MCP Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gptdev
Group=gptdev
WorkingDirectory=/opt/hetzner-dev-workspace
EnvironmentFile=/etc/gpt-dev/gateway.env
Environment=HOME=/var/lib/gpt-dev
Environment=DOCKER_CONFIG=/var/lib/gpt-dev/docker-config
Environment=HANDOFF_INBOX_DIR=/var/lib/gpt-dev/handoffs/incoming
ExecStart=/usr/bin/node apps/mcp-gateway/dist/server.js
Restart=on-failure
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
# Rootless Docker exposes its API below /run/user. Keep home trees read-only
# while allowing the gateway to connect to that user-owned Unix socket.
ProtectHome=read-only
ReadWritePaths=/srv/gpt-hosted /var/lib/gpt-dev
CapabilityBoundingSet=
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
