# infra/systemd/gpt-dev-backup.service
[Unit]
Description=Back up Hetzner Dev Workspace state

[Service]
Type=oneshot
User=gptdev
Group=gptdev
ExecStart=/opt/hetzner-dev-workspace/infra/backup/backup.sh

